# cron/crontab
# XenForo job runner - runs every minute
* * * * * root cd /var/www/html && /usr/local/bin/php cmd.php xf:run-jobs >> /var/log/xenforo-cron.log 2>&1

# Daily backup at 5:00 AM
0 5 * * * root /usr/local/bin/backup.sh >> /var/log/backup.log 2>&1

# Empty line required at end of crontab

